# templates/launchd/com.ian.opennews.rss-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ian.opennews.rss-fetch</string>

  <key>WorkingDirectory</key>
  <string>/Users/ianliu/Documents/Ian_Open_News</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/ianliu/Documents/Ian_Open_News/scripts/local_rss_daily.py</string>
  </array>

  <!-- 評估頻率期間：每小時跑一次（3600 秒）。評估完可改回 StartCalendarInterval 固定時段。 -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/Users/ianliu/Documents/Ian_Open_News/.cache/launchd-rss-fetch.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ianliu/Documents/Ian_Open_News/.cache/launchd-rss-fetch.err</string>
</dict>
</plist>
