# templates/launchd/com.ian.opennews.rss-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ian.opennews.rss-fetch</string>

  <!-- 把下面四個 /ABSOLUTE/PATH/TO/your-repo 換成你 fork 後 repo 的絕對路徑 -->
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/your-repo</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/ABSOLUTE/PATH/TO/your-repo/scripts/local_rss_daily.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>23</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/your-repo/.cache/launchd-rss-fetch.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/your-repo/.cache/launchd-rss-fetch.err</string>
</dict>
</plist>
